# daemon/install/flexitracker.service
[Unit]
Description=flexitracker activity daemon
After=graphical-session.target

[Service]
Type=simple
# Reads ~/.config/flexitracker/config.toml written by `flexitracker configure`.
# Idle detection needs the user's graphical session (X11).
ExecStart=%h/.local/bin/flexitracker
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
