# daemon/install/flexi-worker.service
[Unit]
Description=flexi-worker activity daemon
After=graphical-session.target

[Service]
Type=simple
# Reads ~/.config/flexi-worker/config.toml written by `flexi-worker configure`.
# Idle detection needs the user's graphical session (X11).
ExecStart=%h/.local/bin/flexi-worker
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
